# deploy/suricata-api.service
# Gunicorn API 服務（參考用 systemd unit）
# aaPanel 環境下使用 Python 專案管理器，不需手動使用此檔案。

[Unit]
Description=Suricata Web API (Gunicorn + UvicornWorker)
After=network.target mysql.service redis.service

[Service]
Type=notify
User=www
Group=www
WorkingDirectory=/www/wwwroot/web-suricata-test
EnvironmentFile=/www/wwwroot/web-suricata-test/.env
ExecStart=/www/wwwroot/web-suricata-test/.venv/bin/gunicorn app.main:app \
    -w 4 \
    -k uvicorn.workers.UvicornWorker \
    --bind 127.0.0.1:8000 \
    --access-logfile /www/wwwlogs/suricata-api-access.log \
    --error-logfile /www/wwwlogs/suricata-api-error.log
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
